# launch_agents/com.desktop-reminder.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.desktop-reminder</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/bytedance/Documents/桌面提醒工具/reminder.py</string>
    <string>--schedule</string>
    <string>/Users/bytedance/Documents/桌面提醒工具/schedule.json</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/bytedance/Documents/桌面提醒工具/logs/reminder.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/bytedance/Documents/桌面提醒工具/logs/reminder.err.log</string>
</dict>
</plist>
